# src/etc/supervisor/conf.d/audiofone.conf
[program:audiofone]
command=/opt/futel/src/audiofone.py
user=pi
