# src/etc/supervisor/conf.d/audiofone.conf
[program:audiofone]
command=/opt/futel/src/audiofone.py
user=pi
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/audiofone.log
