# src/etc/supervisor/conf.d/audiofone.conf
[program:audiofone]
command=/opt/futel/src/audiofone.py
user=pi
# gpiozero uses the lgpio backend, which creates a .lgd-nfy* notification file
# in its working directory. supervisord otherwise runs us from / (not writable
# by pi), so point the working dir and HOME/LG_WD at pi's home. Pin the factory
# to lgpio so we fail loudly rather than falling back to RPi.GPIO, whose edge
# detection does not work on recent Pi kernels ("Failed to add edge detection").
directory=/home/pi
environment=HOME="/home/pi",LG_WD="/home/pi",GPIOZERO_PIN_FACTORY="lgpio"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/audiofone.log
